# .config/systemd/user/offlineimap@.timer
[Unit]
Description=OfflineIMAP Service for %i

[Timer]
OnUnitInactiveSec=60s
AccuracySec=1
Unit=offlineimap@%i.service
