# .config/systemd/user/offlineimap.service
[Unit]
Description=OfflineIMAP Service
After=network.target
Conflicts=offlineimap@.service

[Service]
Type=simple
ExecStart=/usr/bin/offlineimap
