# .zsh/profile.d/00-path
typeset -U paths


paths=(
      $HOME/.local/bin
   )

if [[ -d "/usr/games/bin" ]]; then
   paths+="/usr/games/bin"
elif [[ -d "/usr/games" ]]; then
   paths+="/usr/games"
fi

paths+=(
      /usr/bin
      /usr/sbin
      /bin
      /sbin
      $path
   )


path=()
for p in $paths; [[ -d $p ]] && path+=$p
unset p
unset paths


export PATH
