# .zsh/profile.d/10-utils
export VISUAL=nvim
export EDITOR=$VISUAL

export PAGER=less
export MANPAGER=less

if [[ -n "$DISPLAY" ]]; then
    export TERMINAL=germinal
fi
