# zsh/.zsh/profile.d/10-utils.zsh
export VISUAL=nvim
export EDITOR=$VISUAL

export PAGER='bat -p'
export MANPAGER='bat -p'
